# files/docker/env.list
HUSHLINE_USER=hushline
HUSHLINE_GROUP=www-data
APP_PATH=/opt/
APP_FOLDER=hushline
APP_CONF_PATH=/etc/hushline
DOMAIN=staging.hushline.app
EMAIL=hushline@scidsg.org
REDIS_URI=redis://hushline-redis:6379
FLASK_APP=/opt/hushline/hushline/__init__.py
SQLITE_DIR=/var/lib/hushline
SQLALCHEMY_DATABASE_URI=sqlite:////var/lib/hushline/hushline.db
SQLITE_DB_FILE=/var/lib/hushline/hushline.db
SERVICE_FILE=/etc/systemd/system/hushline.service
ONION=your.onion
ENCRYPTION_KEY=
SECRET_KEY=
